# rootfs/etc/services.d/moltbot/run
#!/bin/bash
# OpenClaw gateway service
source /etc/s6-overlay/lib/env-utils.sh
source_env_prefix OPENCLAW_ MOLTBOT_

echo "[moltbot] Starting moltbot gateway..."

export HOME=/home/moltbot
exec s6-envdir /run/s6/container_environment/OPENCLAW_/ \
  s6-envdir /run/s6/container_environment/MOLTBOT_/ \
  s6-setuidgid moltbot bash -l -c 'moltbot gateway --allow-unconfigured'
