# rootfs/etc/services.d/moltbot/run
#!/command/with-contenv bash
# Moltbot gateway service

echo "[moltbot] Starting moltbot gateway..."

# Run as moltbot user with proper environment
export HOME=/home/moltbot
exec s6-setuidgid moltbot bash -c '
  export HOME=/home/moltbot
  export NVM_DIR="$HOME/.nvm"
  [ -s "$NVM_DIR/nvm.sh" ] && . "$NVM_DIR/nvm.sh"
  export PNPM_HOME="$HOME/.local/share/pnpm"
  export PATH="$PNPM_HOME:$PATH"
  exec moltbot gateway --allow-unconfigured
'
