# rootfs/etc/services.d/tailscale/run
#!/command/with-contenv bash

# Run containerboot to configure tailscale for container usage
/usr/local/bin/containerboot
